# fixtures/test.service
[Unit]
Description=example systemd service unit file.

[Service]
ExecStart=echo 1

[Install]
WantedBy=multi-user.target